# core/lib/generators/mno_enterprise/templates/scripts/upstart/app-web-hotrestart.conf
start on <%= app_name %>-hotrestart
task

# Set process owner
setuid app-deployer
setgid app-deployer

# Load App environment 
env RAILS_ENV=<%= environment %>
env APP_HOME=/apps/<%= app_name %>/current
env USER=app-deployer
env HOME=/home/app-deployer

# Load ruby environment
env GEM_HOME=/usr/local/rvm/gems/ruby-2.1.2
env IRBRC=/usr/local/rvm/rubies/ruby-2.1.2/.irbrc
env MY_RUBY_HOME=/usr/local/rvm/rubies/ruby-2.1.2
env PATH=/usr/local/rvm/gems/ruby-2.1.2/bin:/usr/local/rvm/gems/ruby-2.1.2@global/bin:/usr/local/rvm/rubies/ruby-2.1.2/bin:/usr/local/rvm/bin:/sbin:/bin:/usr/sbin:/usr/bin:/opt/aws/bin
env GEM_PATH=/usr/local/rvm/gems/ruby-2.1.2:/usr/local/rvm/gems/ruby-2.1.2@global
env RUBY_VERSION=ruby-2.1.2

script  
  # Reload server
  cd $APP_HOME
  exec bundle exec pumactl -F scripts/<%= environment %>/puma.rb phased-restart
end script