# lib/generators/mno_enterprise/templates/scripts/upstart/app-web.conf
start on starting <%= app_name %>
stop on stopping <%= app_name %>