# lib/generators/mno_enterprise/templates/scripts/upstart/app-web-server.conf
start on starting <%= app_name %>-web
stop on stopping <%= app_name %>-web
respawn

# Set process owner
setuid app-deployer
setgid app-deployer

# Load App environment 
env RAILS_ENV=<%= environment %>
env APP_HOME=/apps/<%= app_name %>/current
env USER=app-deployer
env HOME=/home/app-deployer

# Load ruby environment
env GEM_HOME=/usr/local/rvm/gems/ruby-2.1.2
env IRBRC=/usr/local/rvm/rubies/ruby-2.1.2/.irbrc
env MY_RUBY_HOME=/usr/local/rvm/rubies/ruby-2.1.2
env PATH=/usr/local/rvm/gems/ruby-2.1.2/bin:/usr/local/rvm/gems/ruby-2.1.2@global/bin:/usr/local/rvm/rubies/ruby-2.1.2/bin:/usr/local/rvm/bin:/sbin:/bin:/usr/sbin:/usr/bin:/opt/aws/bin
env GEM_PATH=/usr/local/rvm/gems/ruby-2.1.2:/usr/local/rvm/gems/ruby-2.1.2@global
env RUBY_VERSION=ruby-2.1.2

pre-start script
  # Enforce permissions
  chown -R $USER.$USER $APP_HOME/log/
  chown -R $USER.$USER $APP_HOME/tmp/pids/
  chown -R $USER.$USER $APP_HOME/tmp/cache/
  chown -R $USER.$USER $APP_HOME/tmp/sockets/
end script

script
  cd $APP_HOME
  exec bundle exec puma -C scripts/<%= environment %>/puma.rb
end script