# lib/generators/mno_enterprise/templates/scripts/upstart/app.conf
start on runlevel [2345]
stop on runlevel [016]

pre-start script

bash << "EOF"
  mkdir -p /apps/<%= app_name %>/current/log
  chown -R app-deployer.app-deployer /apps/<%= app_name %>/current/log
EOF

end script